# src/stagepi-package/DEBIAN/postinst
#!/bin/bash
# Stage Pi: Open source stagebox firmware
# Copyright (C) 2025 Bama Box ltd.
#
# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, version 3 of the License.

# This program is distributed in the hope that it will be useful,
#but WITHOUT ANY WARRANTY; without even the implied warranty of
#MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
#GNU General Public License for more details.

# You should have received a copy of the GNU General Public License
# along with this program. If not, see <https://www.gnu.org/licenses/>.

set -e

echo "Post-install script running..."

# Source directory inside the package
SOURCE_DIR="/usr/local/stagepi/installer-files"

# Copy configuration files and directories to /etc
install -d /etc/linuxptp
install -m 644 "$SOURCE_DIR/etc/linuxptp/ptp4l.conf" /etc/linuxptp/
install -d /etc/bluetooth
install -m 644 "$SOURCE_DIR/etc/bluetooth/main.conf" /etc/bluetooth/
install -m 644 "$SOURCE_DIR/etc/shairport-sync.conf" /etc/

# Copy supervisor configuration
install -d /etc/supervisor/conf.d
install -m 644 "$SOURCE_DIR/etc/supervisor/conf.d/stagepi.conf" /etc/supervisor/conf.d/

# Copy systemd service files
install -d /etc/systemd/system
install -m 644 "$SOURCE_DIR/etc/systemd/system/stagepi-ble.service" /etc/systemd/system/
install -m 644 "$SOURCE_DIR/etc/systemd/system/supervisord.service" /etc/systemd/system/
install -m 644 "$SOURCE_DIR/etc/systemd/system/ptp4l.service" /etc/systemd/system/
install -m 644 "$SOURCE_DIR/etc/systemd/system/first-boot.service" /etc/systemd/system/

# Check if we're in a container or if systemd is available
if [ -d /run/systemd/system ] && [ "$(ps -p 1 -o comm=)" = "systemd" ]; then
    # We have systemd, safe to use systemctl
    systemctl daemon-reload
    systemctl enable stagepi.service
    systemctl start stagepi.service
    # Stop and disable old stagepi-ui service if it exists
    systemctl stop stagepi-ui.service 2>/dev/null || true
    systemctl disable stagepi-ui.service 2>/dev/null || true

    # Reload systemd to recognize new/changed service files
    systemctl daemon-reload

    # Enable services
    systemctl enable shairport-sync || true
    systemctl enable bluetooth || true
    systemctl enable supervisord.service || true
    systemctl enable stagepi-ble.service || true
    systemctl enable ptp4l.service || true
    systemctl enable first-boot.service || true

    # Restart services
    systemctl restart shairport-sync || true
    systemctl restart bluetooth || true
    systemctl restart stagepi-ble.service || true
    systemctl restart ptp4l.service || true

    # Start supervisor (which will start the backend)
    systemctl restart supervisord.service || true

    # Reconfigure locales
    #dpkg-reconfigure locales || true
else
    echo "SystemD not available (container environment detected)"
    echo "Services will be managed by Supervisor"
    # Modify supervisor config for container environment (no pi user)
    sed -i 's/^user=pi$/user=root/' /etc/supervisor/conf.d/stagepi.conf
    # Supervisor will be started by the container's entrypoint/CMD
fi

exit 0
