# src/stagepi-package/usr/local/stagepi/installer-files/etc/supervisor/conf.d/stagepi.conf
# Stage Pi: Open source stagebox firmware
# Copyright (C) 2025 Bama Box ltd.
#
# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, version 3 of the License.

# This program is distributed in the hope that it will be useful,
#but WITHOUT ANY WARRANTY; without even the implied warranty of
#MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
# GNU General Public License for more details.

# You should have received a copy of the GNU General Public License
# along with this program. If not, see <https://www.gnu.org/licenses/>.

[program:stagepi-backend]
command=/usr/bin/uvicorn main:app --host 0.0.0.0 --port 8000
directory=/usr/local/stagepi/ui
user=pi
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/stagepi-backend.log
stderr_logfile=/var/log/supervisor/stagepi-backend-error.log
environment=PYTHONUNBUFFERED="1"
